# packaging/systemd/scrobload.service
[Unit]
Description=Scrobload (fetch Last.fm scrobbles and download tracks)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=scrobload
Group=scrobload
EnvironmentFile=-/etc/scrobload.env
ExecStart=/bin/bash -lc '/usr/bin/scrobload $SCROBLOAD_ARGS'

[Install]
WantedBy=multi-user.target
